# files/watchlog.service
[Unit]
Description=My watchlog service
[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/watchlog
ExecStart=/opt/watchlog.sh $WORD $LOG
